# scripts/systemd/forexbot-dashboard.service
[Unit]
Description=ForexBot v2 — Django dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PROJECT_DIR@
ExecStart=@PROJECT_DIR@/venv/bin/python run_dashboard.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Logs: journalctl -u forexbot-dashboard -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
